# supervisord.conf
[supervisord]
nodaemon=true

[program:webmin]
command=/etc/init.d/webmin start

